# applications/luci-app-wifimanager/root/etc/uci-defaults/01_perms
#!/bin/sh
chmod +x /usr/bin/wifimanager
chmod +x /etc/init.d/wifimanager
/etc/init.d/wifimanager enable
/etc/init.d/wifimanager start
return 0
