# applications/luci-app-ftp/root/etc/uci-defaults/01_perms
#!/bin/sh
chmod +x /usr/bin/ftp-server
chmod +x /etc/init.d/ftp
/etc/init.d/ftp enable
return 0
